# templates/mcelog.service
[Service]
ExecStartPre=
ExecStart=/usr/sbin/mcelog --ignorenodev --daemon --foreground
